# docker/mgmt-cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel for management app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/management
ExecStart=%h/management/docker/bin/cloudflared tunnel --no-autoupdate --url http://127.0.0.1:8080
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
